# packaging/scripts/postinst
#!/bin/sh
# postinst: runs after `dpkg -i craftdeck_*.deb` (i.e. `apt install craftdeck`).
# See requirements.md FR-41~45, FR-42a~d, FR-43, FR-43a.
set -e

# 1. Dedicated system account for the craftdeckd binary path/config, even
#    though the service itself runs as root (see craftdeck.service comment
#    and requirements.md FR-43a) -- used for file ownership under
#    /etc/craftdeck, not for running the daemon.
if ! getent passwd craftdeck >/dev/null; then
	adduser --system --group --home /var/lib/craftdeck --no-create-home craftdeck
fi

mkdir -p /var/lib/craftdeck /etc/craftdeck
chown craftdeck:craftdeck /var/lib/craftdeck /etc/craftdeck

# 2. Java 8/17/21/25 (Adoptium Temurin, FR-42/FR-42a/FR-42g) are declared as
#    real apt Depends (see packaging/nfpm.yaml) rather than installed here
#    with a nested `apt-get install` -- this script runs *inside* the outer
#    `apt-get install craftdeck` transaction, which already holds
#    /var/lib/dpkg/lock-frontend, so a second apt-get call from here would
#    deadlock against it instead of just failing (confirmed on real
#    hardware: "Could not get lock ... held by process ... (apt-get)").
#    Registering the Adoptium repo itself is install.sh's job, since that
#    has to happen *before* the outer transaction starts for apt to be able
#    to resolve these Depends at all.

# 3. Enable and start the service so the daemon is running immediately
#    after install, with no extra steps required from the user (FR-43).
systemctl daemon-reload
systemctl enable --now craftdeck.service

# 4. Apply the console log retention/rotation policy (FR-16,
#    /etc/systemd/journald.conf.d/craftdeck.conf) immediately -- journald
#    only picks up config changes on restart, not on next log write.
systemctl restart systemd-journald

# 5. Surface the admin URL right away -- otherwise the operator has no way
#    to know the service is even reachable, let alone at what address,
#    right after `apt install craftdeck` finishes. hostname -I's first
#    entry is the primary LAN IP on a typical single-NIC Raspberry Pi
#    setup; falls back to "<라즈베리파이-IP>" if that lookup ever comes up
#    empty (e.g. no network yet) rather than printing a blank address.
LOCAL_IP=$(hostname -I 2>/dev/null | awk '{print $1}')
echo ""
echo "CraftDeck가 설치되어 실행 중입니다."
echo "관리 페이지 접속: http://${LOCAL_IP:-<라즈베리파이-IP>}:8080"
echo "처음 접속하면 관리자 계정을 만드는 화면이 나옵니다."
echo ""

exit 0
