# packaging/systemd/craftdeck.service
[Unit]
Description=CraftDeck Minecraft server manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
; Runs as root: craftdeckd calls systemd-run to create per-instance
; transient units (with cgroup limits) that run as a dedicated, fixed
; system user per instance -- see requirements.md FR-43a~c,
; ARCHITECTURE.md section 5.1. This User setting is intentionally NOT
; dropped to a service account, since craftdeckd itself needs root to
; create those transient units and per-instance accounts (useradd/chown).
User=root
ExecStart=/usr/bin/craftdeckd
Restart=on-failure
RestartSec=5
Environment=CRAFTDECK_DATA_DIR=/var/lib/craftdeck
Environment=CRAFTDECK_MASTER_KEY=/etc/craftdeck/master.key
Environment=CRAFTDECK_LISTEN_ADDR=:8080

[Install]
WantedBy=multi-user.target
